# backup/cron.d/backup.job
5 1 * * * /usr/bin/backup > /var/log/backup.log 2>&1
